# start-lc.service
[Unit]
Description=LCEC Auto Start Service

[Service]
User=cnc
ExecStart=bash /home/cnc/scripts/start-lc.sh

[Install]
WantedBy=multi-user.target
